# alarmclock_menu.service
[Unit]
Description=menu_handler

[Service]
User=pi
ExecStart=/home/pi/alarmclock-webradio/service_menu.sh

[Install]
WantedBy=multi-user.target
