# alarmclock_rotary.service
[Unit]
Description=rotary_handler

[Service]
User=pi
ExecStart=/home/pi/alarmclock-webradio/service_rotary.sh

[Install]
WantedBy=multi-user.target
