# alarmclock_radio.service
[Unit]
Description=radio_handler

[Service]
User=rsalmon
ExecStart=/app/alarmclock-webradio/service_radio.sh

[Install]
WantedBy=multi-user.target
